# units/redis@.service.d/publish-6379.conf
[Service]
# Expose the Redis port (6379)
Environment='DOCKER_PUBLISH_6379=6379:6379'
